# config/crontab
0 0 * * * python /app/source/manage.py marquee_number
